# packages/gentoo/sys-fs/sanoid/files/sanoid.cron
* * * * * root TZ=UTC /usr/bin/sanoid --cron
